# src/res/linux/systemd.service
[Unit]
Description=espanso daemon

[Service]
ExecStart={{{espanso_path}}} daemon
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target

