# espanso/src/res/linux/systemd.service
[Unit]
Description=espanso

[Service]
ExecStart={{{espanso_path}}} launcher 
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target

